# curator4-Vale.conf
[program:curator4-Vale]
command=python3 /home/curator4/curator4-Vale/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-Vale/test.err.log
stdout_logfile=/home/curator4/curator4-Vale/test.out.log